# rc.local
for ugen in $(usbconfig | awk -F: '{print $1}' | sort) ; do cmd="/usr/sbin/usbconfig -d $ugen power_save"; $($cmd 2> /dev/null) && echo $cmd ; done
